# target/linux/ar71xx/base-files/etc/init.d/bootok
#!/bin/sh /etc/rc.common

START=99

start() {
	. /lib/functions.sh

	local board=$(board_name)
	if [ $board = "dlb-5" ]; then
		# Clear linux fail counter for the active firmware image.
                active_fw=$(fw_printenv -n active)
                fw_setenv linux_fail${active_fw}=0
	fi
}
