# package/my_package/appledns/files/etc/uci-defaults/luci-i18n-appledns-zh-cn
uci set luci.languages.zh_cn='普通话 (Chinese)'; uci commit luci
