# package/linino/linino-scripts/files/etc/init.d/delete_uhttpd_cert
#!/bin/sh /etc/rc.common

START=49

start() {
	CERT=`uci get uhttpd.main.cert`
	KEY=`uci get uhttpd.main.key`
	rm -f $CERT $KEY
}

