# package/linino/linino-scripts/files/etc/init.d/tmpudev-init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=55

start() {
    # trigger udev
    /sbin/udevadm trigger &
}
