# package/linino/linino-scripts/files/etc/init.d/udev-init
#!/bin/sh /etc/rc.common
#Init script for udevd

START=30

start() {
    #----------- START the UDEV daemon --------------
    /sbin/udevd --daemon &
    }

stop() {
    killall udevd
}
