# target/linux/mvebu/base-files/etc/uci-defaults/04_mambafan
#!/bin/sh
#
# Copyright (C) 2017 LEDE-Project.org
#

. /lib/mvebu.sh

board=$(mvebu_board_name)

case "$board" in
armada-xp-linksys-mamba)
        # Set one execution of the fan script during every start (waiting for cron)
        echo "/sbin/fan_ctrl.sh" > /etc/rc.local
        # Set fan script execution in crontab
        echo "# mamba fan script runs every 5 minutes" >> /etc/crontabs/root
        echo "*/5 * * * * /sbin/fan_ctrl.sh" >> /etc/crontabs/root
        # Execute one time after initial flash (instead of waiting 5 min for cron)
        /sbin/fan_ctrl.sh
    ;;
*)
esac

exit 0
